# traefik/rootfs/etc/services.d/traefik/run
#!/usr/bin/with-contenv bashio

OPTIONS_FILE="/data/options.json"
CONFIG_FILE=$(bashio::config 'config_file')

if [ ! -f "${CONFIG_FILE}" ]; then
  bashio::log.fatal "Config file not found at '${CONFIG_FILE}'."
  bashio::log.notice "Example configuration:\n\n$(cat /etc/traefik/traefik.example.yaml)"
  exit 1
fi

bashio::string.replace "$(cat "${CONFIG_FILE}")" '${INGRESS_PORT}' "$(bashio::addon.ingress_port)" > /etc/traefik/traefik.yaml

ENV_VARS=$(cat "${OPTIONS_FILE}" | jq -r \
  '.env_vars[]? | "\(.key)=\(.value)"'
)

if [ -z "${ENV_VARS}" ]; then
  bashio::log.info "No additional environment variables defined."
else
  bashio::log.info "Additional environment variables:\n$(sed 's;=.*;;' <<< $ENV_VARS)"
fi

bashio::log.info "Starting Traefik..."

if [ -z "${ENV_VARS}" ]; then
  exec /usr/local/bin/traefik
else
  env "${ENV_VARS}" /usr/local/bin/traefik
fi
