# traefik/rootfs/etc/services.d/traefik/run
#!/usr/bin/with-contenv bashio

CONFIG_FILE=$(bashio::config 'config_file')

if [ ! -f "${CONFIG_FILE}" ]; then
  bashio::log.fatal "Config file not found at '${CONFIG_FILE}'."
  bashio::log.notice "Example configuration:\n\n$(cat /etc/traefik/traefik.example.yaml)"
  exit 1
fi

bashio::string.replace "$(cat "${CONFIG_FILE}")" '${INGRESS_PORT}' "$(bashio::addon.ingress_port)" > /etc/traefik/traefik.yaml

bashio::log.info "Starting Traefik..."

/usr/local/bin/traefik
